# etc/periodic/15min/ntp
#!/bin/bash

exec ntpd -n -q -d -p ntp.zerotier.com >/dev/shm/ntp.out 2>&1
