# etc/local.d/93edge-services.start
#!/bin/bash

if [ ! -f /edge-debug ]; then
	/opt/edge/run-mgr.sh &
	/opt/edge/run-webui.sh &
fi

exit 0
