# etc/local.d/93start-service-watchers.start
#!/bin/bash

/opt/edge/run-mgr.sh &
/opt/edge/run-webui.sh &

exit 0
