# etc/local.d/99blink.start
#!/bin/bash

/bin/bash -c '/bin/sleep 5; /bin/nice -n -20 /opt/edge/mgr/utils/led-control boot' &

/opt/edge/monitor-network-access.sh >>/dev/null 2>&1 &

exit 0
